# in/init=redhat=wicd.in
#! /bin/sh
#
# /etc/init.d/wicd
#
# chkconfig: - 98 02
# description: Wicd is a wireless and wired network manager for Linux.

# Source function library.
. /etc/rc.d/init.d/functions

# This script was borrowed from the Wicd wiki page on building for Fedora
# http://wicd.net/wiki/doku.php?id=fedora

WICD_BIN=%SBIN%wicd

case "$1" in
        start)
                echo -n "Starting wicd services: "
                daemon $WICD_BIN 2> /dev/null
                touch /var/lock/subsys/wicd
        ;;
        stop)
                echo -n "Shutting down wicd services: "
                killall wicd 2> /dev/null
                rm -f /var/lock/subsys/wicd
        ;;
        status)
                if pidofproc wicd > /dev/null ; then
                        echo "wicd is running."
                else
                        status wicd
                fi
        ;;
        restart | reload)
                $0 stop
                $0 start
        ;;
        *)
                echo "Usage: wicd {start|stop|status|reload|restart}"
                exit 1
        ;;
esac

